# supervisord.conf
[supervisord]
nodaemon=true

[program:deploy-script]
command=bash -c /usr/local/tomcat/bin/deploy-script.sh
